# etc/profile.d/docker.sh
alias dcs='docker ps --format "table {{.Names}}\t{{.Ports}}\t{{.Status}}"'
